# home/Library/LaunchAgents/com.leonardoacosta.nexus-listener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.leonardoacosta.nexus-listener</string>

    <!--
      Bun-based replacement for the bash `nexus-notifier.sh listen`. The
      bash version stalls on macOS 3.2 process-substitution races; this
      uses native fetch streaming so there's no pipe-handling involved.
      See ~/dev/personal/installfest/home/dot_local/share/nexus-listener.ts for the source.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/bun</string>
        <string>run</string>
        <string>/Users/leonardoacosta/.local/share/nexus-listener.ts</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/leonardoacosta</string>
        <!-- nexus-agent runs on homelab, not this Mac -- reach it via Tailscale
             MagicDNS (stable across IP re-registration/churn). REVERTED
             2026-07-20 (if-2vm1 stopgap): the AdGuard-authoritative
             leonardoacosta.dev name broke real tooling (paste-image.sh,
             rotate-keys.sh) -- this Mac's live DNS config has zero
             resolver rules for leonardoacosta.dev at all right now, so
             nothing routes there (not an AdGuard outage -- AdGuard itself
             answers fine when queried directly; the split-DNS push from
             hl-oz6 just isn't reaching this client). Back to MagicDNS
             until that's diagnosed with GUI/admin-console access. The old
             100.91.88.16 was a stale/wrong IP that never matched a live
             device -- this FQDN is the one that actually works. -->
        <key>NEXUS_URL</key>
        <string>http://homelab.tail296462.ts.net:7400</string>
        <!-- Pre-bootstrap defaults — bootstrapSettings() in the listener
             overrides these with values from /notifications/settings on
             startup. Match the bash plist defaults for parity. -->
        <key>TTS_ENABLED</key>
        <string>1</string>
        <key>BANNER_ENABLED</key>
        <string>1</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/leonardoacosta/.local/state/nexus-listener.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/leonardoacosta/.local/state/nexus-listener.stderr.log</string>

    <!-- ThrottleInterval prevents launchd from respawn-storming if the
         listener exits immediately (e.g., bun missing or syntax error). -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
