# home/Library/LaunchAgents/com.leonardoacosta.mic-priority.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.leonardoacosta.mic-priority</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/Users/leonardoacosta/dev/personal/installfest/scripts/mic-priority.sh</string>
    </array>

    <!-- Run at login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Also run every 30 seconds to catch device changes -->
    <key>StartInterval</key>
    <integer>30</integer>

    <!-- Don't keep it running, just execute and exit -->
    <key>KeepAlive</key>
    <false/>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/mic-priority.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/mic-priority.error.log</string>
</dict>
</plist>
