# home/Library/LaunchAgents/com.leonardoacosta.cloudpc-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.leonardoacosta.cloudpc-tunnel</string>

    <!-- ControlMaster=no + ControlPath=none override ~/.ssh/config's global
         `Host *` ControlMaster auto / ControlPersist 10m (added for
         paste-image.sh handshake reuse). Without this, any interactive
         `ssh cloudpc` from a Terminal session creates a background
         multiplexed master that THIS launchd job's ssh hands the -D 1080
         forwarding off to, then exits 0 -- launchd sees a clean exit (not
         a crash) so KeepAlive doesn't relaunch it, and the tunnel dies
         silently once that interactive session's 10-minute ControlPersist
         timer expires. Found + fixed on the homelab systemd twin of this
         unit, 2026-07-20, diagnosing "SOCKS is down". -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/ssh</string>
        <string>-D</string>
        <string>1080</string>
        <string>-N</string>
        <string>-o</string>
        <string>ControlMaster=no</string>
        <string>-o</string>
        <string>ControlPath=none</string>
        <string>-o</string>
        <string>ServerAliveInterval=30</string>
        <string>-o</string>
        <string>ServerAliveCountMax=3</string>
        <string>-o</string>
        <string>ExitOnForwardFailure=yes</string>
        <string>cloudpc</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>NetworkState</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/sbin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/leonardoacosta</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/leonardoacosta/.local/logs/cloudpc-tunnel.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/leonardoacosta/.local/logs/cloudpc-tunnel.err.log</string>
</dict>
</plist>
