# home/Library/LaunchAgents/com.leonardoacosta.validate-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.leonardoacosta.validate-proxy</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/leonardoacosta/dev/personal/installfest/scripts/validate-proxy.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Check every 60 seconds -->
    <key>StartInterval</key>
    <integer>60</integer>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/sbin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/leonardoacosta</string>
        <key>DOTFILES</key>
        <string>/Users/leonardoacosta/dev/personal/installfest</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/leonardoacosta/.local/logs/validate-proxy.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/leonardoacosta/.local/logs/validate-proxy.err.log</string>
</dict>
</plist>
