# homelab/usb-boot/systemd/homelab-bootstrap.service
[Unit]
Description=Homelab Bootstrap Service
Documentation=https://github.com/leonardoacosta/Installfest
After=network-online.target docker.service
Wants=network-online.target
ConditionPathExists=!/var/lib/homelab-bootstrap-complete

[Service]
Type=oneshot
User=%i
ExecStartPre=/bin/sleep 10
ExecStart=/bin/bash /home/%i/homelab-bootstrap.sh
ExecStartPost=/bin/touch /var/lib/homelab-bootstrap-complete
StandardOutput=append:/var/log/homelab-bootstrap.log
StandardError=append:/var/log/homelab-bootstrap.log
TimeoutStartSec=3600
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
