# homeserver/homeserver.service
[Unit]
Description=Homeserver Podman Stack
After=network-online.target podman.socket
Wants=network-online.target
RequiresMountsFor=%t/containers

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/%u/homeserver
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Start services
ExecStart=/usr/bin/podman-compose up -d

# Stop services
ExecStop=/usr/bin/podman-compose down

# Reload services
ExecReload=/usr/bin/podman-compose restart

# For rootless podman
User=%u
Group=%u

# Restart policy
Restart=on-failure
RestartSec=30

[Install]
WantedBy=default.target
