# homelab/systemd/restore-tailscale-routes.service
[Unit]
Description=Restore Tailscale IP Routes and iptables Rules
Documentation=https://github.com/yourusername/homelab/docs/TAILSCALE-ADGUARD-IPTABLES-ISSUE.md
After=docker.service tailscaled.service network-online.target
Wants=network-online.target
# Don't require tailscaled.service since Tailscale might be in Docker
ConditionPathExists=/usr/local/bin/restore-tailscale-routes.sh

[Service]
Type=oneshot
ExecStart=/usr/local/bin/restore-tailscale-routes.sh
RemainOnExit=yes
StandardOutput=journal
StandardError=journal

# Restart policy
Restart=on-failure
RestartSec=10s

# Security
PrivateTmp=yes
NoNewPrivileges=no
# Note: NoNewPrivileges must be 'no' because we need root for iptables

# Environment
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="TS_ROUTES=172.20.0.0/16,172.21.0.0/16"

[Install]
WantedBy=multi-user.target
