# home/dot_config/systemd/user/ssh-agent.service.d/homelab-keys.conf
# Auto-load the GitHub-registered key into the socket-activated ssh-agent on
# every start (agent identities are memory-only and don't survive an agent
# restart/reboot otherwise). `-` tolerates a missing key file so this drop-in
# stays harmless on a machine that never generated this key. See if-kjbg —
# without an agent+key, bd's git+ssh Dolt remote fails all push/pull with
# "Permission denied (publickey)".
[Service]
# ExecStartPost runs as a fresh process — it does NOT inherit the socket path
# ssh-agent is listening on, so ssh-add needs it set explicitly or it exits 2
# ("could not connect to authentication agent") every time (confirmed live).
Environment=SSH_AUTH_SOCK=%t/ssh-agent.socket
ExecStartPost=-/usr/bin/ssh-add %h/.ssh/id_ed25519_cc_20260721
