# launchd/com.leonardoacosta.cmux-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.leonardoacosta.cmux-bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/leonardoacosta/.claude/scripts/bin/cmux-bridge</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/leonardoacosta/.claude/logs/cmux-bridge.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/leonardoacosta/.claude/logs/cmux-bridge.err.log</string>
</dict>
</plist>
