# home/dot_config/systemd/user/dotenv-key-rotation-check.timer
[Unit]
Description=Run dotenv-key-rotation-check.service daily (if-jhqf)

[Timer]
OnCalendar=daily
RandomizedDelaySec=15m
Persistent=true
Unit=dotenv-key-rotation-check.service

[Install]
WantedBy=timers.target
