# home/dot_config/systemd/user/file-server.service
[Unit]
Description=File server for cmux embedded browser
Documentation=https://github.com/nyaptor/if
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/dev/personal/installfest/scripts/file-server.py --port 8787 --bind 0.0.0.0
Restart=on-failure
RestartSec=5
Environment=FILE_SERVER_PORT=8787
Environment=FILE_SERVER_BIND=0.0.0.0

[Install]
WantedBy=default.target
