# home/dot_config/systemd/user/mesh-heartbeat.service
[Unit]
Description=Mesh + broker heartbeat (Tailscale mac/cloudpc, mx-broker /health, SOCKS tunnel)
Documentation=file:%h/dev/personal/installfest/openspec/changes/add-auth-resilience-and-hygiene/proposal.md

[Service]
Type=oneshot
ExecStart=/usr/bin/bash %h/dev/personal/installfest/scripts/mesh-heartbeat.sh
EnvironmentFile=-%h/.env
StandardOutput=append:%h/.local/state/mesh-heartbeat.systemd.log
StandardError=append:%h/.local/state/mesh-heartbeat.systemd.log
Restart=no
