# home/dot_config/systemd/user/nexus-dashboard.service
[Unit]
Description=Nexus Dashboard — Next.js web UI for session management
After=network-online.target
Wants=network-online.target
# Stop auto-restart after 5 failures in 30 seconds; prevents runaway crash loops.
StartLimitBurst=5
StartLimitIntervalSec=30

[Service]
Type=simple
WorkingDirectory=%h/dev/nx/apps/nextjs

# Require a production build before starting.
# Run: cd ~/dev/nx/apps/nextjs && pnpm build
ExecStartPre=/bin/sh -c 'test -d "%h/dev/nx/apps/nextjs/.next" || { echo "ERROR: .next/ not found — run pnpm build in apps/nextjs first"; exit 1; }'
# Port 3100 is opened persistently in /etc/ufw/before.rules — no runtime iptables needed.
ExecStart=%h/.local/share/pnpm/pnpm start

Restart=always
RestartSec=5

Environment=NODE_ENV=production
Environment=PORT=3100
Environment=HOSTNAME=0.0.0.0
Environment=PATH=%h/.local/share/mise/shims:%h/.local/share/pnpm:%h/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin

# Load machine-local secrets: NEXUS_AGENTS, POSTGRES_URL, etc.
EnvironmentFile=-%h/.env

# Security hardening (mirrors nexus-agent.service)
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/dev/nx/apps/nextjs/.next
ReadWritePaths=%h/dev/nx/apps/nextjs/.next/cache

# Resource limits
LimitNOFILE=65536
MemoryMax=1G
TasksMax=512

[Install]
WantedBy=default.target
