# home/dot_config/systemd/user/ropen.service
[Unit]
Description=ropen — long-lived multi-mount HTTP file server for browsing dev folder over Tailscale
Documentation=https://github.com/nyaptor/if
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Resolve the live Tailscale IPv4 address fresh on every start (handles address
# changes across reboots/reconnects) and write it to a runtime-dir EnvironmentFile
# — ropen-server.py defaults ROPEN_BIND to 127.0.0.1 (harden-mesh-file-servers),
# so this is what makes cross-machine ropen/gopen/sopen/mopen/iopen actually
# reach the server again. Falls back to 127.0.0.1 (loopback-only, same as the
# script's own default) if tailscale is down/unreachable at start — fails safe,
# never re-widens to 0.0.0.0.
ExecStartPre=/bin/sh -c 'echo "ROPEN_BIND=$(tailscale ip -4 2>/dev/null || echo 127.0.0.1)" > %t/ropen-bind.env'
EnvironmentFile=%t/ropen-bind.env
ExecStart=%h/.local/bin/ropen --serve
Restart=always
RestartSec=5
# Port can be overridden via systemctl --user edit ropen
Environment=ROPEN_PORT=8889
# State dir lives under /tmp/ropen-<uid>; recreated on startup.
# Logs land in journalctl --user -u ropen (stdout/stderr captured by systemd).

[Install]
WantedBy=default.target
